# deb/odoohost-1.0.0/etc/profile.d/01-environment.sh
#!/bin/bash
export PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/root/bin"
export PGPASSWORD="toBe72/*"
export PGUSER=odoo
export PGDATABASE=processcontrol
export PGHOST=localhost
export TZ=Europe/Madrid
export LANG C.UTF-8