# deb/odoohost-1.0.0/etc/systemd/system/jupyterlab.service
[Unit]
Description=Jupyterlab
After=network.target

[Service]
Environment=PATH=/opt/odoo.sh/jupyterlab/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Type=simple
SyslogIdentifier=jupyterlab
PermissionsStartOnly=true
User=odoo
Group=odoo
ExecStart=jupyter-lab --ip=0.0.0.0 --ServerApp.base_url=/odoo-sh/jupyterlabr
KillMode=mixed
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target