# deb/odoohost-1.0.0/etc/systemd/system/odoo.service
[Unit]
Description=Odoo Open Source ERP and CRM
After=network.target

[Service]
Type=simple
SyslogIdentifier=odoo
PermissionsStartOnly=true
User=odoo
Group=odoo
ExecStart=/home/odoo/venv/python3 /home/odoo/src/odoo/odoo-bin
KillMode=mixed
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target