# deb/odoohost-1.0.0/etc/systemd/system/odoo.service
[Unit]
Description=Odoo Open Source ERP and CRM
After=network.target

[Service]
Restart=always
Type=simple
SyslogIdentifier=odoo
PermissionsStartOnly=true
User=odoo
Group=odoo
ExecStart=/root/bin/start-init
KillMode=mixed
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target