# deb/odoohost-1.0.0/etc/systemd/system/odoo.service
[Unit]
Description=Odoo Open Source ERP and CRM
After=network.target

[Service]
EnvironmentFile=/etc/environment
Type=simple
SyslogIdentifier=odoo
PermissionsStartOnly=true
User=odoo
Group=odoo
ExecStartPre=extract_addons_path
ExecStart=/home/odoo/venv/bin/python3 /home/odoo/src/odoo/odoo-bin --database=${PGDATABASE} --logfile=/home/odoo/logs/odoo.log --addons-path=${ADDONS_PATH}
KillMode=mixed
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target